# scripts/s6/act_runner/run
#!/usr/bin/env bash

set -e

exec run.sh
