# scripts/s6/docker/run
#!/usr/bin/env bash

set -e

exec dockerd-entrypoint.sh
